# recipes-python/haslate/files/haslate.service.in
[Unit]
Description=haslate startup script
After=xserver-nodm.service

[Service]
ExecStart=/home/@USER@/haslate/haslate.py
Restart=always
WorkingDirectory=~
User=@USER@
Environment="DISPLAY=:0.0"

[Install]
WantedBy=multi-user.target
